# sysutils/jobd/pkg-plist
bin/jobctl
lib/libjob.so.0
lib/libjob.so
lib/libjob.a
/etc/rc.d/jobd
sbin/jobd
%%MANPREFIX%%/man/man1/jobctl.1.gz
%%MANPREFIX%%/man/man1/jobcfg.1.gz
%%MANPREFIX%%/man/man5/job.5.gz
%%MANPREFIX%%/man/man8/jobd.8.gz
@dir /etc/job.d
@dir /var/cache/jobd
@dir /var/db/jobd
