# sysutils/powerdxx/pkg-plist
bin/loadplay
bin/loadrec
/etc/rc.d/powerdxx
lib/libloadplay.so
%%MANPREFIX%%/man/man1/loadplay.1.gz
%%MANPREFIX%%/man/man1/loadrec.1.gz
%%MANPREFIX%%/man/man8/powerd++.8.gz
%%MANPREFIX%%/man/man8/powerdxx.8.gz
sbin/powerd++
sbin/powerdxx
%%PORTDOCS%%%%DOCSDIR%%/README.md
%%PORTDOCS%%%%DOCSDIR%%/refman.pdf
