# sysutils/qchroot/pkg-plist
bin/qchroot
/etc/rc.d/qchroot.bootime
%%MANPREFIX%%/man/man8/qchroot.8.gz
