# net-mgmt/bpft/pkg-plist
bin/trafdump
bin/traflog
bin/trafsave
bin/trafstart
bin/trafstat
bin/trafstop
/etc/rc.d/trafd.sh
/etc/traflog.format
%%MANPREFIX%%/man/man1/traflog.1.gz
%%MANPREFIX%%/man/man1/trafstat.1.gz
%%MANPREFIX%%/man/man8/trafd.8.gz
%%MANPREFIX%%/man/man8/trafstatd.8.gz
sbin/trafd
sbin/trafstatd
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/WhatsNew
@dir var/bpft
